# dash-nginx-uwsgi-docker-master/supervisord.conf
[program:uwsgi]
environment=PATH='/opt/conda/envs/crispritz/bin'
command=/opt/conda/envs/crispritz/bin/uwsgi --ini /etc/uwsgi/uwsgi.ini --die-on-term --need-app